# gunicorn_flask.conf
[program:gunicorn]
command=/var/www/competition/flask/bin/gunicorn -w 4 -b 127.0.0.1:8000 app:app 
directory=/var/www/competition
autostart=true
autorestart=true
redirect_stderr=true

